# eoip-watch.service
[Unit]
Description=EOIP Watch

[Service]
Type=oneshot
ExecStart=/opt/xtain-hcloud-network/watch
